# ansible-docker-ubuntu/roles/services/files/dockerclean.timer
[Unit]
Descrition=run docker clean up

[Timer]
OnBootSec=11min
OnUnitActiveSec=7d

[Install]
WantedBy=multi-user.target
