# ansible-docker-ubuntu/roles/services/files/dockerclean.timer
[Unit]
Description=run docker clean up

[Timer]
OnCalendar=weekly
Persistent=true
OnBootSec=60s

[Install]
WantedBy=multi-user.target
