# ansible-docker-ubuntu/roles/services/files/trashcli.service
[Unit]
Description=empty trash daily older 3

[Service]
User=regular_user
Group=regular_user
Type=oneshot
ExecStart=/usr/bin/trash-empty 3
